# Apple_Homekit/HomeKit ADK 2.0/HomeKit ADK 2.0/ADK/SDK-Samples/POSIX/Raspi/Tools/helpers/Pi/lib/systemd/system/acme-lightbulb.service
[Unit]
Description=Acme LightbulbLED server on port 8080

[Service]
WorkingDirectory=/home/pi/acme-lightbulb/bin
ExecStart=/home/pi/acme-lightbulb/bin/LightbulbLED
Restart=on-failure

[Install]
WantedBy=multi-user.target
